# crontab
*/5 * * * * gcalcli --nocolor agenda "`date`" > /tmp/gcalcli_agenda.txt
